# supervisord.conf
[supervisord]
nodaemon=true

[program:titan]
command=/opt/titan-server-0.4.2/bin/titan.sh start
